# packages/kanno-core/root/etc/init.d/kanno
#!/bin/sh /etc/rc.common
# KannoProxy init script

START=99
STOP=10

PROG=/usr/bin/kanno

start() {
    local enabled
    enabled=$(uci -q get kanno.global.enabled 2>/dev/null)
    [ "${enabled:-1}" = "0" ] && return 0
    $PROG start
}

stop() {
    $PROG stop
}

reload() {
    $PROG reload
}
